# root/etc/init.d/warp-cron
#!/bin/sh /etc/rc.common
# WARP Cron Job for China IP update
# Copyright (C) 2024

START=99

boot() {
    # 姣忓懆鏃ュ噷鏅?鐐规洿鏂颁腑鍥絀P鍒楄〃
    echo "0 3 * * 0 /usr/bin/warp-update-china >/dev/null 2>&1" >> /etc/crontabs/root
    /etc/init.d/cron restart
}